# templates/sysv
#!/bin/bash

# This file was autogenerated by node-deb {{ node_deb_version }}

### BEGIN INIT INFO
# Provides:          {{ node_deb_package_name }}
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: {{ node_deb_package_name }}
# Description:       {{ node_deb_package_description }}
### END INIT INFO

DAEMON_DIR={{ node_deb_install_path }}/{{ node_deb_package_name }}
DAEMON_USER={{ node_deb_user }}

SU=/bin/su

. /lib/lsb/init-functions

do_start () {
    log_daemon_msg "Starting {{ node_deb_package_name }} service"
    $SU -l $DAEMON_USER --shell=/bin/bash -c "$DAEMON_DIR/bin/{{ node_deb_executable_name }}" || return 2
    log_end_msg $?
}
do_stop () {
    log_daemon_msg "Stopping {{ node_deb_package_name }} service"
    $SU -l $DAEMON_USER --shell=/bin/bash -c "$DAEMON_DIR/bin/{{ node_deb_shutdown_name }}" || return 2
    log_end_msg $?
}

case "$1" in

    start|stop)
	do_${1}
	;;

    restart|reload|force-reload)
	do_stop
	do_start
	;;

    status)
	status_of_proc "$DAEMON_NAME" "$DAEMON" && exit 0 || exit $?
	;;

    *)
	echo "Usage: /etc/init.d/$DAEMON_NAME {start|stop|restart|status}"
	exit 1
	;;

esac
exit 0
